# meraki-alerts.service
[Unit]
Description=Meraki Alerts Service
After=network.target

[Service]
Type=simple
User=<your-user>
WorkingDirectory=/opt/merakialerts
ExecStart=/usr/bin/python3 /opt/merakialerts/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
